# systemd/careeros-nextjs.service
[Unit]
Description=CareerOS Next.js Frontend
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/job-tracker
Environment="PATH=/usr/bin:/bin"
Environment="NODE_ENV=production"
Environment="PORT=3000"
EnvironmentFile=/home/YOUR_USERNAME/job-tracker/.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Performance
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
