# systemd/careeros-ollama.service
[Unit]
Description=Ollama Service for CareerOS
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
Environment="OLLAMA_HOST=0.0.0.0:11434"
Environment="OLLAMA_MODELS=/home/YOUR_USERNAME/.ollama/models"
Environment="OLLAMA_MAX_LOADED_MODELS=2"
Environment="OLLAMA_NUM_PARALLEL=2"
ExecStart=/usr/local/bin/ollama serve
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Memory limit for 8GB server (allocate 6GB to Ollama)
MemoryMax=6G
MemoryHigh=5.5G

[Install]
WantedBy=multi-user.target
