# systemd/careeros-flask.service
[Unit]
Description=CareerOS Flask Backend
After=network.target careeros-ollama.service redis.service
Requires=careeros-ollama.service redis.service

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/job-tracker/services/ai-backend
Environment="PATH=/home/YOUR_USERNAME/job-tracker/services/ai-backend/venv/bin"
EnvironmentFile=/home/YOUR_USERNAME/job-tracker/services/ai-backend/.env
ExecStart=/home/YOUR_USERNAME/job-tracker/services/ai-backend/venv/bin/gunicorn --bind 0.0.0.0:5000 --workers 2 --timeout 120 --access-logfile /var/log/careeros/flask-access.log --error-logfile /var/log/careeros/flask-error.log wsgi:app
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Performance
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
